# im-backend.service
# Systemd 服务配置文件 - IM 后端服务
# 部署位置: /etc/systemd/system/im-backend.service
# 
# 使用方法:
#   sudo systemctl daemon-reload          # 重新加载配置
#   sudo systemctl start im-backend       # 启动服务
#   sudo systemctl stop im-backend        # 停止服务
#   sudo systemctl restart im-backend     # 重启服务
#   sudo systemctl status im-backend      # 查看状态
#   sudo systemctl enable im-backend      # 开机自启
#   sudo systemctl disable im-backend     # 禁用开机自启

[Unit]
Description=IM Backend Service
After=network.target mysql.service redis.service elasticsearch.service
Wants=mysql.service redis.service

[Service]
Type=simple
User=root
Group=root

# 工作目录（jar 包所在目录）
WorkingDirectory=/opt/im/backend

# 启动命令
ExecStart=/usr/bin/java \
    -Xms512m \
    -Xmx1024m \
    -XX:+UseG1GC \
    -XX:MaxGCPauseMillis=200 \
    -Dspring.profiles.active=prod \
    -jar /opt/im/backend/im-0.0.1-SNAPSHOT.jar

# 停止命令
ExecStop=/bin/kill -15 $MAINPID

# 重启策略
Restart=on-failure
RestartSec=10s

# 超时设置
TimeoutStartSec=60
TimeoutStopSec=30

# 日志输出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=im-backend

# 安全设置
PrivateTmp=true
NoNewPrivileges=true

# 环境变量（可选，也可以在 application-prod.properties 中配置）
# Environment="MYSQL_URL=jdbc:mysql://localhost:3306/im"
# Environment="MYSQL_USERNAME=root"
# Environment="MYSQL_PASSWORD=your_password"
# Environment="REDIS_HOST=localhost"
# Environment="REDIS_PORT=6379"

[Install]
WantedBy=multi-user.target

